# config/ruflo-sync.service
[Unit]
Description=Sync and build ruflo from upstream
After=network-online.target
Wants=network-online.target
OnFailure=ruflo-notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/claude/src/ruflo-patch
EnvironmentFile=/home/claude/.config/ruflo/secrets.env
Environment=HOME=/home/claude
Environment=PATH=/home/claude/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/claude/src/ruflo-patch/scripts/sync-and-build.sh
TimeoutStartSec=3600
MemoryMax=32G
CPUQuota=800%
StandardOutput=journal
StandardError=journal
