# config/ruflo-notify-failure@.service
[Unit]
Description=Notify on %i failure

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo "[$(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ)] ALERT: %i failed" >> /tmp/ruflo-pipeline-failures.log'
